# launchagent/dev.stack3.job-hunt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.stack3.job-hunt</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/stack3/.local/bin/claude</string>
        <string>-p</string>
        <string>/job-hunt</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/stack3/.claude/logs/job-hunt.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/stack3/.claude/logs/job-hunt-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/stack3/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/stack3</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
